# cron.weekly/zfs_autoscrub
#!/bin/bash

for pool in $(zpool list -H -o name)
do
        zpool scrub $pool
done
