# cron.weekly/zfs_autoscrub
#!/bin/bash

for pool in $(zpool list -H -o name)
do
    health=$(zpool list -H -o health $pool)
    if [ "$health" = "ONLINE" ] ; then
        zpool scrub $pool
    fi
done
